# docker/app/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=php /var/www/html/artisan queue:work --tries=3
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/supervisord_queue_work.log

[program:crontab]
command=php /var/www/html/artisan schedule:work
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/supervisord_schedule_work.log
